# ci/docker/rootfs/etc/services.d/laravel-horizon/run
#!/usr/bin/with-contenv bash

cd /app
exec s6-setuidgid nginx php artisan horizon
